# server/deploy/supervisord_conf.ini
[program:ext_comment_backend]
user=root
command=gunicorn -c deploy/prd_comment.py run:app
autostart=true
autorestart=true
stopasgroup=true
directory=/var/www/ext_comment/server